# sv/elogin/run
#!/bin/bash
exec &>/var/log/all/.log
sv start /etc/service/dbus &>/dev/null || exit
exec /lib64/elogind/elogind
