# sv/acpid/run
#!/bin/sh
exec 2>&1
[[ -S /dev/log ]] || exit
exec acpid -c /etc/acpi/events -f
