# sv/getty-generic/finish
#!/bin/sh
tty=${PWD##*-}
exec utmpset =q $tty
